# dist/crontab/crontab
*/10 * * * *      /usr/bin/fetch -o /dev/null http://localhost/of/index.php?module=aggrtraff
*/5 * * * *      /usr/bin/fetch -o /dev/null http://localhost/of/index.php?module=aggrgdata
