# dist/crontab/crontab.preconf
*/10 * * * *      /usr/bin/fetch -o /dev/null "http://127.0.0.1/of/?module=aggrtraff"
*/5 * * * *      /usr/bin/fetch -o /dev/null "http://127.0.0.1/of/?module=aggrgdata"
30 3 * * *       /usr/bin/fetch -o /dev/null "http://127.0.0.1/of/?module=rotator"

