# dist/crontab/crontab.preconf
*/10 * * * *      /usr/bin/fetch -o /dev/null http://localhost/of/?module=aggrtraff
*/5 * * * *      /usr/bin/fetch -o /dev/null http://localhost/of/?module=aggrgdata

